# debian/webblock/lib/systemd/system/webblock.service
 [Unit]
Description=webblock-service Başlatıldı.
After=multi-user.target

[Service]
Type=simple
#Environment=export DISPLAY=:0
#ExecStartPre=/usr/bin/bash 
ExecStart=/usr/bin/webblock-service
#ExecStop=/usr/share/webblock-service
ExecStop=/usr/share/webblock/script/webblock-clean
Restart=on-failure

[Install]
WantedBy=multi-user.target
